# src/sys/hhl7.socket
[Unit]
Description=HHL7 incoming socket

[Socket]
ListenStream=5376
MaxConnections=64

[Install]
WantedBy=sockets.target
