# src/sys/hhl7.socket
[Unit]
Description=HHL7 incoming socket

[Socket]
ListenStream=192.168.0.28:5376
MaxConnections=64

[Install]
WantedBy=sockets.target
